# src/init/https_gateway
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=95
STOP=10

PROG=/usr/sbin/https-gateway
CONF=https_gateway
GENERATED_CONF=/etc/nginx/conf.d/https-gateway.conf

start_service() {
	local enabled
	config_load "$CONF"
	config_get enabled global enabled '0'

	[ "$enabled" = "1" ] || return 0

	# Generate nginx conf, sync ACME, configure dnsmasq
	"$PROG" apply
}

stop_service() {
	# Remove generated nginx conf
	rm -f "$GENERATED_CONF"

	# Restore uhttpd to default ports if needed
	local cur_http
	cur_http=$(uci -q get uhttpd.main.listen_http)
	if echo "$cur_http" | grep -q '127.0.0.1:8080'; then
		uci set uhttpd.main.listen_http='0.0.0.0:80'
		uci set uhttpd.main.listen_https='0.0.0.0:443'
		uci commit uhttpd
		/etc/init.d/uhttpd restart
	fi

	# Reload nginx to drop our server blocks
	/etc/init.d/nginx reload 2>/dev/null
}

reload_service() {
	start_service
}

service_triggers() {
	procd_add_reload_trigger "$CONF"
}
